# etc/salt/supervisord-boxed.conf
[supervisord]
nodaemon=true

[program:salt]
command=salt-minion

[program:nginx]
command=nginx -g "daemon off;"
